# rpi/hipache/supervisord/hipache.ini
[program:hipache]
command=/usr/bin/hipache -c /etc/hipache/config.json
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true